# bin/node-daemon.conf
# /etc/init/node-daemon.conf
start on started network
stop on stopping network
respawn
expect daemon
exec https-proxy-daemon